# etc/systemd/bloop.service
[Unit]
Description=Bloop Scala build server

[Service]
ExecStart=%h/.bloop/bloop server

[Install]
WantedBy=default.target
